# rc.dnsbl
#!/bin/sh
#
# $OpenBSD: rc.template,v 1.5 2011/03/10 13:46:59 ajacoutot Exp $

daemon="/usr/local/sbin/dnsbl-divert"
#daemon_flags=""
#daemon_user="root"

. /etc/rc.d/rc.subr

rc_cmd $1
